# cookbooks/thin/templates/default/thin_init.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          <%= @name %>_<%= @env_name %>
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
# Short-Description: <%= @name %> initscript
# Description:       <%= @name %> thin controller
### END INIT INFO
# Original author: Forrest Robertson

DAEMON=/usr/bin/thin
SCRIPT_NAME=/etc/init.d/<%= @name %>_<%= @env_name %>
CONFIG_FILE="<%= node.default["thin"]["fs_etc_root"] %>/<%= @name %>.yml"
case "$1" in
  start)
        if [[ "$USER" == "ubuntu" ]]
				then
					$DAEMON start -C $CONFIG_FILE
				else
					su ubuntu -c $DAEMON start -C $CONFIG_FILE
				fi
        ;;
  stop)
        $DAEMON stop -C $CONFIG_FILE
        ;;
  restart)
        $DAEMON restart -C $CONFIG_FILE
        ;;
  *)
        echo "Usage: $SCRIPT_NAME {start|stop|restart}" >&2
        exit 3
        ;;
esac

:
